# supervisor/demo_project.conf
[program:millstone]
command=/home/ubuntu/.local/bin/gunicorn_django -b 127.0.0.1:8000
directory=/home/ubuntu/demo_project/demo_project
autostart=1
autorestart=True
redirect_stderr=True
user=ubuntu
